# packaging/aur/vpn-tray.install
post_install() {
  cat <<'EOF'
==> vpn-tray requires a sudoers rule for privileged helper actions.
==> Install and customize the example file, then validate it with visudo:
==>   sudo install -m 0440 /usr/share/doc/vpn-tray/examples/sudoers.vpn-tray /etc/sudoers.d/vpn-tray
==>   sudo visudo -cf /etc/sudoers.d/vpn-tray
==> Enable user service:
==>   systemctl --user enable --now vpn-tray.service
EOF
}

post_upgrade() {
  post_install
}
